# recipes-support/ntp/files/ntpd.service
[Unit]
Description=Network Time Service
After=network.target

[Service]
ExecStart=/usr/bin/ntpd -g -n

[Install]
WantedBy=multi-user.target
